# test/support/docker/etc/profile.d/aliases.sh
alias la='ls -1A'
alias ll='ls -al'
alias l='ls -1'
